# DEBIAN/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable azurefile-dockervolumedriver
systemctl start azurefile-dockervolumedriver
